# templates/vault.upstart.erb
# vault Agent (Upstart unit)
###########################################################################################################
# this file has been put in place by the jsok/vault Puppet module (https://forge.puppetlabs.com/jsok/vault)
# any changes will be overwritten if Puppet is run again
###########################################################################################################
description "vault server"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [06]

env VAULT=<%= scope['vault::bin_dir'] %>/vault
env CONFIG=<%= scope['vault::config_dir'] %>/config.json
env USER=<%= scope['vault::user'] %>
env GROUP=<%= scope['vault::group'] %>
env PID_FILE=/var/run/vault.pid
env LOG_FILE=/var/log/vault.log

script
    export GOMAXPROCS=${GOMAXPROCS:-<%= scope['vault::num_procs'] %>}
    [ -e /etc/default/$UPSTART_JOB ] && . /etc/default/$UPSTART_JOB
    exec >> $LOG_FILE 2>&1
    exec start-stop-daemon -u $USER -g $GROUP -p $PID_FILE -x $VAULT -S -- <%= scope['vault::service_type'] %> -config=$CONFIG <%= scope['vault::service_options'] %>
end script

respawn
respawn limit 10 10
kill timeout 10
